# accounts/modules/addons/ms365backup/deploy/systemd/ms365-worker-fleet.timer
[Unit]
Description=MS365 worker fleet cron timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
Unit=ms365-worker-fleet.service

[Install]
WantedBy=timers.target
