# accounts/modules/addons/ms365backup/deploy/systemd/ms365-worker-release-sync.timer
[Unit]
Description=Run MS365 worker release pull sync every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Unit=ms365-worker-release-sync.service

[Install]
WantedBy=timers.target
